# restoreboot.service
[Unit]
Description=Restore Boot configuration

[Service]
ExecStart=/usr/bin/restoreboot

[Install]
WantedBy=multi-user.target
