# restoreboot.service
[Unit]
Description=Restore Boot configuration

[Service]
ExecStart=/opt/justbrowsing/restoreboot

[Install]
WantedBy=multi-user.target
